# systemd/bafu.service
[Unit]
Description=Unit to download hydroweb.xml from BAFU

[Service]
User=www-data
Group=www-data
Type=oneshot
ExecStart=/srv/swisshydroapi/bafu.sh
